# apod.cron

0 8 * * * ~/go/src/github.com/timmyw/apod/getandsetapod.sh
